# ansible/roles/consul/templates/consul.openrc
#!/sbin/openrc-run

export CHECKPOINT_DISABLE=1
name="$SVCNAME"
command="/usr/local/bin/$SVCNAME"
command_args="agent -config-dir=/usr/local/etc/consul"
command_background=true
command_user="consul:consul"

pidfile="/var/run/$SVCNAME.pid"

depend() {
	need net
}
